# scripts/auto_rk/rc.local
#!/bin/sh
# 开机自启脚本（量产模式）- RK3562 版本
# 仅做系统级初始化，业务逻辑全部交给app_start.sh


# 等待USB摄像头枚举完成（RK3562 启动较快，可缩短等待时间）
sleep 5


# 设置系统时间（兜底）
date -s "2026-06-19 00:00:00"

# 启动业务系统（后台运行）
/mnt/nfs/run_on_board_rk3562/auto_rk/app_start.sh &

exit 0