# scripts/auto/rc.local
#!/bin/sh
# 开机自启脚本（量产模式）
# 仅做系统级初始化，业务逻辑全部交给app_start.sh

# 【关键】全局导出环境变量（sysinit阶段必备）
export APP_HOME="/root/run_on_board"
export LD_LIBRARY_PATH=$APP_HOME/openh264:$APP_HOME/libjpeg:$APP_HOME/mnn:$APP_HOME/libyuv

# 等待USB摄像头枚举完成（最稳9秒，可根据实际情况调整为5-10秒）
sleep 9

# 挂载SD卡
mount /dev/mmcblk0p1 /mnt/sdcard 2>/dev/null

# 设置系统时间（兜底）
date -s "2026-05-29 00:00:00"

# 启动业务系统（后台运行）
/root/run_on_board/auto/app_start.sh &

exit 0